# images/logs/config/supervisord.conf
[supervisord]
nodaemon=true

[program:rsyslogd]
command=/bin/bash -c "/usr/sbin/rsyslogd -n"

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"
